# Library/LaunchAgents/io.nickolaskraus.system-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <!-- Runs `system-update` daily at 08:00. -->
  <dict>
    <key>Label</key>
    <string>io.nickolaskraus.system-update</string>
    <!--
      launchd does not expand environment variables in plist values, so an
      absolute path here would pin this file to a single machine. It does set
      HOME for a gui/<uid> agent, so the job runs via `sh -c` and lets the
      shell expand $HOME. The redirect stands in for StandardOutPath and
      StandardErrorPath, which are subject to the same no-expansion rule.
      `exec` replaces the shell, so ProcessType and LowPriorityIO still apply
      to system-update itself. `install` creates the log directory; the
      redirect fails if it does not exist.
      -->
    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-c</string>
      <string>exec "$HOME/bin/system-update" &gt;&gt; "$HOME/.local/state/system-update.log" 2&gt;&amp;1</string>
    </array>
    <!--
      `system-update` runs `dotfiles install`, and that script reloads this
      agent. Booting out a job kills its process tree, so an unguarded reload
      would terminate the very run that triggered it, and every step below the
      reload would be skipped. This marker is inherited by the whole tree and
      lets `install` recognize that it is running under this agent and leave
      the reload alone.
      -->
    <key>EnvironmentVariables</key>
    <dict>
      <key>SYSTEM_UPDATE_AGENT</key>
      <string>1</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>8</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <key>ProcessType</key>
    <string>Background</string>
    <key>LowPriorityIO</key>
    <true/>
  </dict>
</plist>
